# Sunshine/v0.7.7/sunshine.service
[Unit]
Description=Sunshine Gamestream Server for Moonlight

[Service]
ExecStart=sunshine

[Install]
WantedBy=graphical-session.target
